# lenovo-agent/qubes-kvm-agent.service
# TEMPLATE — setup-lenovo.sh generates the real service dynamically.
# These placeholders show the structure; actual paths are filled at install time.
[Unit]
Description=Qubes KVM Development Agent
Documentation=https://github.com/QubesOS/qubes-core-admin
After=network-online.target libvirtd.service
Wants=network-online.target

[Service]
Type=simple
User=%USER%
Group=%USER%
WorkingDirectory=%PROJECT_DIR%
ExecStart=%AGENT_DIR%/.venv/bin/python -m uvicorn agent:app --host 0.0.0.0 --port 8420 --log-level info --app-dir %AGENT_DIR%
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

Environment=HOME=%HOME%
Environment=PATH=%AGENT_DIR%/.venv/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
